# docker/app/files/supervisord.conf
[program:logstash]
command=/logstash-1.4.2/bin/logstash -f /etc/files/logstash.conf

[program:shortener]
command=uwsgi --ini /www/shortener/conf/uwsgi.ini
